# sys/openkongqi.cron
# /etc/cron.d/openkongqi - okq crontab
# number of days to keep cache files
CACHE_PATH=/srv/openkongqi/_cache/
CACHE_TTL=7
# vacuum the cache for files older than CACHE_TTL days every day at 1:00
0 1 * * *     www-data    find ${CACHE_PATH} -type f -mtime +${CACHE_TTL} -delete
